# slackandcode.conf
description "SlackAndCode yesod app"
start on runlevel [2345];
stop on runlevel [!2345];
respawn
chdir /code
exec /code/dist/build/scaffold/scaffold
